# systemd/wizsmith_onboarding.service
[Unit]
Description=WizSmith Onboarding Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/wizsmith_edge_hub/services/onboarding/main.py
WorkingDirectory=/opt/wizsmith_edge_hub
Restart=always
EnvironmentFile=/opt/wizsmith_edge_hub/.env

[Install]
WantedBy=multi-user.target
