# systemd/wizsmith_frigate_updater.service
[Unit]
Description=WizSmith Edge Hub - Frigate Updater
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/wizsmith_edge_hub/services/frigate_updater/main.py
WorkingDirectory=/opt/wizsmith_edge_hub
EnvironmentFile=/opt/wizsmith_edge_hub/.env
Restart=always

[Install]
WantedBy=multi-user.target
